# initscripts/cluevpn.service
[Unit]
Description=ClueVPN

[Service]
Type=simple
ExecStart=/usr/sbin/cluevpn -c /etc/cluevpn -f

[Install]
WantedBy=multi-user.target
